# scandir_supervisor.conf
[program:scan_dir]
command=python scan_dir.py                                                                                              ;
directory=/data/scan_dir                                                                                                ;
user=work                                                                                                               ;
startsecs=0                                                                                                             ;
stopwaitsecs=0                                                                                                          ;
autostart=true                                                                                                          ;
autorestart=true                                                                                                        ;
stdout_logfile=/data/scan_dir/log/gunicorn.log                                                                          ;
stderr_logfile=/data/scan_dir/log/gunicorn.err                                                                          ;
redirect_stderr=true                                                                                                    ;
stdout_logfile_maxbytes=20MB                                                                                            ;
stdout_logfile_backups=20                                                                                               ;                                                                             ;